# odoo/supervisord.conf.d/bash_supervisord.conf
[supervisord]
nodaemon=true

[program:odoo]
user=odoo
environment=HOME="/opt/odoo"
redirect_stderr=true
command=/bin/bash
