# odoo/supervisord.conf.d/nginx_supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
user=root
stdout_logfile=/var/log/nginx/nginx_supervisord.txt
redirect_stderr=true
command=/bin/bash -c "nginx"
