# odoo/supervisor.conf
[supervisord]
nodaemon=true

[program:postresql]
command=/usr/lib/postgresql/9.3/bin/postgres --config_file=/etc/postgresql/9.3/main/postgresql.conf -D /var/lib/postgresql/9.3/main

[program:odoo]
command=/opt/odoo/odoo/odoo.py -c /etc/odoo/openerp-server.conf
